# launchd/com.user.markdown-reviewer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.markdown-reviewer</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT__/bin/server</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PORT</key>
        <string>__PORT__</string>
        <key>REVIEW_ROOT</key>
        <string>__REVIEW_ROOT__</string>
        <!-- markdown-reviewer doesn't actually use the user DB (boilerplate
             leftover); :memory: avoids needing a writable cwd at launch. -->
        <key>DATABASE_DSN</key>
        <string>:memory:</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/markdown-reviewer/markdown-reviewer.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/markdown-reviewer/markdown-reviewer.err.log</string>
</dict>
</plist>
